# rootfs/etc/services.d/carbonoz/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Carbonoz service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Carbonoz SolarAutopilot..."

exec /usr/bin/carbonoz.sh